# modules/codesearch/files/supervisor/cs-server.conf
[program:cs-server]
command=/home/nelhage/sw/bin/node /home/nelhage/codesearch/web/cs_server.js
directory=/home/nelhage/codesearch
user=nelhage
autostart=true
autorestart=true
