# modules/codesearch/files/supervisor/web-server.conf
[program:web-server]
command=/home/nelhage/sw/bin/node /home/nelhage/codesearch/web/web_server.js
directory=/home/nelhage/codesearch
user=nelhage
autostart=true
autorestart=true
