# modules/codesearch/files/supervisor/web-server.conf
[program:web-server]
command=node /home/nelhage/codesearch/web/web_server.js --production
directory=/home/nelhage/codesearch
user=nelhage
autostart=true
autorestart=true
